# deploy/jobfitr-harvest.service
# jobfitr harvest — the demoted BASELINE inflow. Builds the broad jobs.json snapshot,
# so the pool has a floor even for searches nobody has run live yet. Live per-search
# fetches (on /api/score) own freshness now; this just keeps a broad baseline warm.
# Oneshot; fired by jobfitr-harvest.timer.
# Install: /etc/systemd/system/jobfitr-harvest.service (done by bootstrap.sh).
#
# Under blue-green the harvest runs from whichever slot is PRODUCTION — harvest-active.sh
# resolves /etc/jobfitr/active-slot at run time, because a unit file can't. It writes only
# the SHARED /opt/jobfitr/data/jobs.json; each slot's own jobs.db then pulls that snapshot
# in via store.sync_snapshot() (the background sync in the web process). One harvest, no
# doubled job-API traffic, and every slot stays current without a restart. The --watchlist
# (the ~94-company ATS depth lane) comes from the same slot, so harvester and list always
# ship together — see harvest-active.sh.

[Unit]
Description=jobfitr harvest (build the jobs.json snapshot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=jobfitr
Group=jobfitr
EnvironmentFile=/etc/jobfitr/jobfitr.env
ExecStart=/opt/jobfitr/bin/harvest-active.sh

# The harvest was ~8 min against 94 companies. With the resolved universe (several
# hundred boards) plus Workday's per-job description fetches it is expected to run
# 25-40 min. Bound it: a single hung upstream fetch must not leave the job running
# into tomorrow and block the next night's harvest. 90 min is generous headroom over
# the expected worst case, and a kill here is recoverable — jobs.json is written
# atomically, so a timeout leaves yesterday's snapshot intact rather than a partial one.
TimeoutStartSec=90min

# Hardening — may write only the shared snapshot directory. (The harvest's best-effort
# store upsert lands in the shared jobs.db named by jobfitr.env, which lives here too.)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/jobfitr/data
